# deploy/supervisor/conf.d/celery-worker.ini
[program:celery-worker]
user = root
directory = /app
command = celery --app pelican_publisher worker  -l info --concurrency 1

redirect_stderr = true
stdout_logfile = /log/celery-worker.log
