# deploy/supervisor/conf.d/django-tasks-worker.ini
[program:django-tasks-worker]
user = root
directory = /app
command = ./manage.py db_worker

redirect_stderr = true
stdout_logfile = /log/django-tasks-worker.log