# deploy/supervisor/conf.d/django-toosimple-q.ini
[program:django-toosimple-q]
user = root
directory = /app
command = ./manage.py worker

redirect_stderr = true
stdout_logfile = /log/django-toosimple-q.log