# deploy/supervisor/celery-worker.conf
[program:celery-worker]
user = root
directory = /app/pelican_publisher
command = celery --app pelican_publisher worker  -l info --concurrency 1

redirect_stderr = true
stdout_logfile = /pp-data/celery-worker.log
