# deploy/supervisor/celery-worker.conf
[program:celery-worker]
user = root
directory = /app/pelican_publisher
command = celery -A pelican_publisher worker  -l info

redirect_stderr = true
stdout_logfile = /var/log/celery-worker.log
